# deploy/systemd/aredn-ldap-bridge.service
[Unit]
Description=AREDN LDAP Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=aredn-ldap-bridge
Group=aredn-ldap-bridge
WorkingDirectory=/opt/aredn-ldap-bridge
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/aredn-ldap-bridge/src
ExecStart=/opt/aredn-ldap-bridge/venv/bin/python -m aredn_ldap_bridge
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=3
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
